# etc_example/systemd.system.fix-intel-power.service
[Unit]
Description=Fix Intel CPU power permission

[Service]
ExecStart=/usr/local/sbin/fix-intel-power.sh

[Install]
WantedBy=default.target
